# systemd/alpaca-worker.service
[Unit]
Description=Alpaca Worker - File Processing Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=alpaca
WorkingDirectory=/home/alpaca/alpaca
Environment="PATH=/home/alpaca/alpaca/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FILEWATCHER_API_URL=http://localhost:8081"
ExecStart=/home/alpaca/alpaca/venv/bin/python /home/alpaca/alpaca/worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
