# scripts/alpaca.service
[Unit]
Description=ALPACA RAG - Prefect Worker
After=network.target postgresql.service

[Service]
Type=simple
User=alpaca
WorkingDirectory=/home/alpaca/alpaca
Environment="PATH=/home/alpaca/alpaca/venv/bin"
ExecStart=/home/alpaca/alpaca/venv/bin/python /home/alpaca/alpaca/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
